# docker/bigpods/experience-pod/supervisor/supervisord.prod.conf

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
silent=false
minfds=4096
minprocs=500
strip_ansi=true
identifier=experience-pod

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=admin
password=dreamscape_experience_123

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=admin
password=dreamscape_experience_123

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ===============================================
# NGINX - Static Assets & VR Content Server
# Priority 10 (starts first)
# ===============================================
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=SIGQUIT
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
user=root
environment=NGINX_WORKER_PROCESSES="auto"

# ===============================================
# Panorama Service - VR Content Management
# Priority 20 (starts after NGINX)
# ===============================================
[program:panorama-service]
command=node server.js
directory=/app/panorama
user=nodejs
autostart=true
autorestart=true
priority=20
startsecs=15
startretries=3
stopwaitsecs=10
stopsignal=SIGTERM
stdout_logfile=/var/log/supervisor/panorama-service.log
stderr_logfile=/var/log/supervisor/panorama-service-error.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=5
environment=NODE_ENV=production,PORT=3006,SERVICE_NAME="panorama-service",VR_CONTENT_PATH="/usr/share/nginx/html/vr",NODE_OPTIONS="--max-old-space-size=512"

# ===============================================
# Gateway Service - API Gateway & Proxy
# Priority 30 (starts after panorama)
# ===============================================
[program:gateway-service]
command=node server.js
directory=/app/gateway
user=nodejs
autostart=true
autorestart=true
priority=30
startsecs=15
startretries=3
stopwaitsecs=10
stopsignal=SIGTERM
stdout_logfile=/var/log/supervisor/gateway-service.log
stderr_logfile=/var/log/supervisor/gateway-service-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
environment=NODE_ENV=production,PORT=3000,SERVICE_NAME="gateway-service",NODE_OPTIONS="--max-old-space-size=256"

# ===============================================
# VR Content Optimizer - Background Processing
# Priority 40 (starts after services)
# ===============================================
[program:vr-optimizer]
command=python3 /app/scripts/vr_optimizer.py
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=3
stopwaitsecs=15
stopsignal=SIGTERM
stdout_logfile=/var/log/supervisor/vr-optimizer.log
stderr_logfile=/var/log/supervisor/vr-optimizer-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
user=nodejs
environment=VR_SOURCE_PATH="/usr/share/nginx/html/vr",VR_CACHE_PATH="/var/cache/nginx/vr",OPTIMIZATION_QUALITY="85"

# ===============================================
# Experience Pod Health Monitor
# Priority 50 (starts last)
# ===============================================
[program:health-monitor]
command=python3 /app/scripts/experience_health_monitor.py
autostart=true
autorestart=true
priority=50
startsecs=5
startretries=3
stopwaitsecs=5
stopsignal=SIGTERM
stdout_logfile=/var/log/supervisor/health-monitor.log
stderr_logfile=/var/log/supervisor/health-monitor-error.log
stdout_logfile_maxbytes=25MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=25MB
stderr_logfile_backups=3
user=nodejs
environment=MONITOR_INTERVAL="30",ALERT_THRESHOLDS="cpu:80,memory:85,disk:90"

# ===============================================
# Process Groups for Coordinated Management
# ===============================================

# Frontend services group
[group:frontend-services]
programs=nginx
priority=100

# Backend services group  
[group:backend-services]
programs=panorama-service,gateway-service
priority=200

# Optimization and monitoring group
[group:optimization-services]
programs=vr-optimizer,health-monitor
priority=300

# ===============================================
# Event Listeners for Experience Pod Monitoring
# ===============================================

# Crash notification system
[eventlistener:crash-notifier]
command=/app/scripts/experience_crash_notifier.py
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED
buffer_size=100
stdout_logfile=/var/log/supervisor/crash-notifier.log
stderr_logfile=/var/log/supervisor/crash-notifier-error.log
autostart=true
autorestart=unexpected

# Performance monitoring
[eventlistener:performance-monitor]
command=/app/scripts/experience_performance_monitor.py
events=TICK_60
buffer_size=100
stdout_logfile=/var/log/supervisor/performance-monitor.log
stderr_logfile=/var/log/supervisor/performance-monitor-error.log
autostart=true
autorestart=unexpected

# VR content monitoring
[eventlistener:vr-content-monitor]
command=/app/scripts/vr_content_monitor.py
events=TICK_60
buffer_size=50
stdout_logfile=/var/log/supervisor/vr-content-monitor.log
stderr_logfile=/var/log/supervisor/vr-content-monitor-error.log
autostart=true
autorestart=unexpected

# ===============================================
# Resource Management and Limits
# ===============================================

# Memory thresholds for services
[program:memory-guard]
command=/app/scripts/memory_guard.py --threshold=85 --critical=95
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/supervisor/memory-guard.log
stderr_logfile=/var/log/supervisor/memory-guard-error.log
user=nodejs

# ===============================================
# Development and Debugging (commented for production)
# ===============================================

# Hot reload for development (disabled in production)
# [program:hot-reload]
# command=npm run dev
# directory=/app/web
# autostart=false
# autorestart=true
# priority=100
# user=nodejs
# environment=NODE_ENV=development,HMR_PORT=3008

# VR content development server (disabled in production)
# [program:vr-dev-server]
# command=python3 -m http.server 8080
# directory=/usr/share/nginx/html/vr
# autostart=false
# autorestart=true
# priority=100
# user=nodejs